# userpatches/overlay/etc/systemd/system/myservice.service
Description=myservice
 
[Service]
ExecStart=/root/myservice.sh
 
[Install]
WantedBy=default.target
